# files/root/etc/systemd/system/uupd.service.d/only-one-upgrade-at-a-time.conf
[Service]
ExecCondition=/usr/bin/bash -c '! systemctl is-active flatpak-system-update.service'
ExecCondition=/usr/bin/bash -c '! systemctl is-active rpm-ostreed-automatic.service'
