# files/root/etc/systemd/system/flatpak-system-update.service.d/only-one-upgrade-at-a-time.conf
[Service]
ExecCondition=/usr/bin/bash -c '! systemctl is-active rpm-ostreed-automatic.service'
ExecCondition=/usr/bin/bash -c '! systemctl is-active uupd.service'
