# files/root/etc/systemd/system/flatpak-system-update.service.d/only-one-upgrade-at-a-time.conf
[Service]
ExecCondition=/usr/bin/bash -c '! (ps aux | grep -E "rpm-ostree.*upgrade" | grep -vq "grep")'
ExecCondition=/usr/bin/bash -c '! (ps aux | grep -E "uupd" | grep -vq "grep")'
