# files/root/etc/systemd/system/rpm-ostreed-automatic.service.d/only-one-upgrade-at-a-time.conf
[Service]
ExecCondition=/usr/bin/bash -c '! systemctl is-active flatpak-system-update.service'
